# shared/rc.local
#!/bin/sh -e

# TODO: See overlay filesystem option in raspi-config:
# https://github.com/RPi-Distro/raspi-config/commit/5d7664812bbeb8a31cb77b70326fab1b257e0946
# Source: http://ala-paavola.fi/jaakko/doku.php?id=pi-readonly
# Note: https://gist.github.com/gbaman/975e2db164b3ca2b51ae11e45e8fd40a
# Note: https://gist.github.com/thenaterhood/56301c3bf7d784a33ab6c46778aa8179
# Note: https://medium.com/swlh/make-your-raspberry-pi-file-system-read-only-raspbian-buster-c558694de79

df -h

echo "*** Set Timezeone ***"
timedatectl set-timezone Australia/Melbourne

echo "*** Install extra debian packages ***"
apt install -y \
    autofs \
    avrdude \
    avr-libc \
    gpiod \
    pigpio \
    libopenblas-dev \
    liblapack-dev \
    libgmp3-dev \
    libmpfr-dev

echo "*** Install /usr/local packages ***"
mount -t vfat /dev/sdb1 /mnt
tar -C /usr/local -xzf /mnt/*.tgz
ln -s /usr/local/julia-*/bin/julia /usr/local/bin/julia
sudo -u pi julia -e 'using Pkg; Pkg.update()'
cp /mnt/*toml /home/pi/.julia/environments/v1.5/
sudo -u pi julia -e 'using Pkg; Pkg.instantiate()'

echo "*** Install busybox-syslogd ***"
apt install busybox-syslogd -y

echo "*** Disable PrivateTmp ***"
sed -i 's/PrivateTmp=yes/PrivateTmp=no/' /lib/systemd/system/systemd-hostnamed.service

echo "*** Disable HDMI output ***"
echo "hdmi_blanking=2" >> /boot/config.txt

echo "*** Disable SD Card Resize ***"
sed -i 's/init=[^ ]*//' /boot/cmdline.txt

echo "*** Disable Serial Console ***"
sudo raspi-config nonint do_serial 1

echo "*** Enable USB OTG Ethernet ***"
echo "dtoverlay=dwc2" >> /boot/config.txt
sed -i 's/rootwait/rootwait noswap modules-load=dwc2,g_ether g_ether.host_addr=f3-91-83-bb-b6-af g_ether.dev_addr=52-81-f6-87-a5-64/' /boot/cmdline.txt

echo "*** Enable i2c ***"
raspi-config nonint do_i2c 0

echo "*** Configure WIFI ***"
rm /etc/wpa_supplicant/wpa_supplicant.conf
if [ -f /mnt/wpa_supplicant.conf.secret ]
then
    cp /mnt/wpa_supplicant.conf.secret /boot/wpa/wpa_supplicant.conf
else
    cp /mnt/wpa_supplicant.conf /boot/wpa
fi
ln -s /boot/wpa /etc/wpa_supplicant/wpa_supplicant.conf

echo "*** Configure Hostname ***"
cp /etc/hostname /etc/hosts /boot
rm /etc/hostname /etc/hosts
ln -s /boot/hostname /etc/hostname
ln -s /boot/hosts /etc/hosts
cp /mnt/fixhostname.service /etc/systemd/system/
systemctl enable fixhostname.service

echo "*** Enable resolv.conf ***"
cp /etc/resolv.conf /var/run/resolv.conf
rm /etc/resolv.conf
ln -sf /var/run/resolv.conf /etc/resolv.conf
sed -i 's/resolv_conf=.*/resolv_conf=\/var\/run\/resolv.conf/' /etc/resolvconf.conf

echo "*** Enable SSH ***"
systemctl enable ssh
mkdir -p /home/pi/.ssh
echo "ssh-rsa AAAAB3NzaC1yc2EAAAADAQABAAABgQDOic/GsgE3ji9WbNwYAtveGonaPpjA2viCHmjnDOP8HOc/BSFuEkQ4h8dF/wmoG27gj9SwzwGgapQRRHMUTAfzz6aMsisjH1LLR5/AQBUqk/YhyHTfyXxFnM7mOgSTvRQgL/K1dOY+4StHaDR3fKP10Dw4bYqMvksR8QyTv5wnLbJ//wetxv1t/vxh5EhauA8dNQI8uqKngwKmdgmiTxViLDcy1yyJANf3JWZ2i/x+Re9BOfWksK1c4EbjPs582i0xOZ7AwqGiQCj7Z2mZtdjUTVB+dRmTLO2XiQeVj09D7TVDJgaRIHGUGSuvebOwpIcboLvCvTtuMn3XDA5czOjiWS7LNxlarB0/iI4BuQjF416lIS45ctBYXxvYNp4PKp7x6y+ERm1LbNrQBYhrQLo0OwbT09GEXRXoXZ2AbZU5+6pVApE1l146kXVY/u34LMSTC5qwMLU4ZpWbU3g38ylMUbo8bKU+U+ej571t6CH2goVsKtvS+bnHsB8EjgOH70E= pi@raspberrypi.local" > /home/pi/.ssh/authorized_keys
chown -R pi:pi /home/pi/.ssh
chmod 700 /home/pi/.ssh
chmod 600 /home/pi/.ssh/authorized_keys

echo "*** CREATE RAM FILESYSTEM ***"
echo 'none    /tmp            tmpfs   size=28M,mode=01777 0 0' >> /etc/fstab
echo 'none    /var/log         tmpfs   size=16M 0 0' >> /etc/fstab
mount -a

echo "*** DISABLE EXTRA SERVICES ***"
systemctl disable dphys-swapfile.service
systemctl disable apt-daily.service apt-daily.timer apt-daily-upgrade.timer apt-daily-upgrade.service
systemctl disable logrotate.service logrotate.timer
systemctl disable alsa-utils triggerhappy

echo "*** CREATE MOUNT SCRIPTS ***"
echo "alias _rw='sudo mount / -o remount,rw;  sudo mount /boot -o remount,rw'" >> /home/pi/.bashrc
echo "alias _ro='sudo mount / -o remount,ro; sudo mount /boot -o remount,ro'" >> /home/pi/.bashrc

echo "*** MAKE FILESYSTEM READ-ONLY ***"
sed -i 's/vfat    defaults/vfat    defaults,ro/' /etc/fstab
sed -i 's/ext4    defaults/ext4    defaults,ro/' /etc/fstab

rm /etc/rc.local
ln -s /boot/rc.local /etc/rc.local

systemctl reboot

exit 0
